# docs/archive/legacy-capacitor/linux-systemd-copies/logueos-companion.service
[Unit]
Description=LogueOS Companion (SvelteKit local-model chat, port 18769)
Documentation=https://github.com/Dreighto/LogueOS-Companion
# Depends on Ollama (the local model), NOT the kernel gateway — the companion
# must boot and run with the LogueOS kernel completely OFF.
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=dreighto
WorkingDirectory=/home/dreighto/dev/LogueOS-Companion
# OWN env file (NOT the canonical Orchestrator .env) — full kernel independence
# and no EnvironmentFile-vs-Environment precedence trap. Contains the companion
# run-mode + private data paths + provider keys (copied from canonical, mode 600,
# gitignored). To rotate keys: regenerate this file from the canonical .env.
EnvironmentFile=/home/dreighto/dev/LogueOS-Companion/.env
ExecStart=/usr/bin/node build/index.js
KillMode=mixed
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
